# deploy/supervisor/invoicing-app-worker.conf
; Queue worker for the Invoicing App, supervised so it restarts on crash and
; survives reboots. Every invoice/receipt email, every SMS, recurring-invoice
; generation, and the overdue-invoice sweep are dispatched onto the queue —
; nothing in any of those happens without this running.
;
; Install:
;   sudo cp deploy/supervisor/invoicing-app-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start invoicing-app-worker:*
;
; After any deploy that changes job/service code, restart workers so they
; pick up the new code (a running worker holds old code in memory
; otherwise): `php artisan queue:restart` (deploy/scripts/deploy.sh already
; does this for you).

[program:invoicing-app-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/invoicing-app/artisan queue:work --sleep=3 --tries=3 --timeout=120
directory=/var/www/invoicing-app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/invoicing-app/storage/logs/worker.log
stopwaitsecs=3600

; One worker process is plenty at this app's scale (single business, not a
; multi-tenant SaaS) — raise numprocs only if delivery volume genuinely
; backs up (check `php artisan queue:monitor` or the jobs table depth).
